# ConfigTinker/hotspotawake.service
[Unit]
Description=Effectue une operation iw toutes les 1s.

[Service]
User=MOI
Group=MOI
Restart=always
RestartSec=1s
ExecStart=/home/MOI/bin/hotspotawake.sh

[Install]
WantedBy=multi-user.target

