# ConfigTinker/x11vnc.service
[Unit]
Description=Start x11vnc at startup.
After=multi-user.target
 
[Service]
Type=simple
ExecStart=/usr/bin/x11vnc -forever -loop -noxdamage -noxrecord -repeat -rfbauth /etc/x11vnc.pass -rfbport 5900 -shared -xkb
 
[Install]
WantedBy=multi-user.target
